# roon-web-config.service
[Unit]
Description=Roon Display Web Configuration Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=dietpi
Group=dietpi
WorkingDirectory=/home/dietpi/roon-full-art-display
Environment=PATH=/home/dietpi/roon-full-art-display/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/dietpi/roon-full-art-display/bin/python -m roon_display.web_config --port 8080 --host 0.0.0.0
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=roon-web-config

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/dietpi/roon-full-art-display

[Install]
WantedBy=multi-user.target