# space-cleaner.timer
[Unit]
Description=Run Space Cleaner Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
